# docker/crontab
0 8-19 * * * /app/docker/sweep.sh >> /mnt/logs/cron.log 2>&1
30 8-19 * * * /app/docker/reconcile.sh >> /mnt/logs/cron.log 2>&1
