# docker/crontab
0 8-19 * * * /app/docker/run-me-cron.sh >> /mnt/logs/cron.log 2>&1
